# flask-sample/flaskapp.conf
[program:flaskapp]
directory=/var/www/flaskapp
command=gunicorn -c /var/www/flaskapp/gunicorncfg.py webapp:app